# shared/files/services/obs.service
[Unit]
Description=OBS

[Service]
User=stream
Restart=always
ExecStart=DISPLAY=:1 obs

[Install]
WantedBy=multi-user.target
