# shared/files/services/obs.service
[Unit]
Description=OBS

[Service]
User=stream
Restart=on-failure
ExecStart=obs
Environment="DISPLAY=:1"

[Install]
kantedBy=multi-user.target
