# shared/files/services/pulse.service
[Unit]
Description=Pulseaudio

[Service]
User=stream
Restart=always
ExecStart=pulseaudio -D

[Install]
WantedBy=multi-user.target
